# crontab
0 9 * * * bundle exec rake reminder:send_emails